# daily-queue.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the BullMQ queue worker. Keeps the daily-pipeline worker
  running 24/7, auto-restarting on crash, with caffeinate to block idle sleep.

  Install (one-time):
    cp daily-queue.worker.plist ~/Library/LaunchAgents/com.thanapol.ai-capital.worker.plist
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.thanapol.ai-capital.worker.plist

  Stop / uninstall:
    launchctl bootout gui/$(id -u)/com.thanapol.ai-capital.worker

  Status / PID:
    launchctl list com.thanapol.ai-capital.worker

  Manual restart (after editing this file):
    launchctl bootout gui/$(id -u)/com.thanapol.ai-capital.worker
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.thanapol.ai-capital.worker.plist
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.thanapol.ai-capital.worker</string>

	<!-- caffeinate -i blocks idle sleep; npx-tsx resolves the worker entry.
	     Absolute paths so we don't depend on launchd's shell PATH. -->
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/caffeinate</string>
		<string>-i</string>
		<string>/opt/homebrew/bin/npx</string>
		<string>tsx</string>
		<string>packages/queue/bin/worker.ts</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/thanapold/Desktop/Projects</string>

	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<!-- ThrottleInterval: wait 30s before respawning a freshly-crashed worker
	     so a tight crash loop doesn't burn CPU. -->
	<key>ThrottleInterval</key>
	<integer>30</integer>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>PIPELINE_RUNS_DB</key>
		<string>/Users/thanapold/Desktop/Projects/data/pipeline-runs.db</string>
		<key>DATA_ROOT</key>
		<string>/Users/thanapold/Desktop/Projects/apps</string>
		<key>REDIS_URL</key>
		<string>redis://localhost:6379</string>
		<!-- @common/db's usePostgres() flips the structured + vector stores to
		     Postgres+pgvector when DATABASE_URL is set. Without it,
		     capital-ingestion silently writes embeddings to the local 33GB
		     corrupt LanceDB dir, which is what was hanging the pipeline. -->
		<key>DATABASE_URL</key>
		<string>postgres://thanapold@localhost:5432/ai_capital</string>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/thanapold/Desktop/Projects/logs/queue-worker.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/thanapold/Desktop/Projects/logs/queue-worker.err.log</string>
</dict>
</plist>
